# scripts/systemd/hl-carry-bot.service
[Unit]
Description=Hyperliquid Carry Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/hl-carry-bot
ExecStart=/opt/hl-carry-bot/hl-carry-bot -config /opt/hl-carry-bot/config.yaml
Restart=on-failure
RestartSec=5
User=hlbot
Group=hlbot

[Install]
WantedBy=multi-user.target
